# dns/etc/supervisor/conf.d/dnsdist.conf
[program:dnsdist]
command=/usr/bin/dnsdist
#autostart=yes
#autorestart=yes
#user=root
